# clipboard2path.service
[Unit]
Description=clipboard2path-wsl - Clipboard image to file path converter
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.cargo/bin/clipboard2path-wsl
Restart=on-failure
RestartSec=5
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_RUNTIME_DIR=/run/user/%U

[Install]
WantedBy=default.target
